# src/templates/license-manager-agent.timer
[Unit]
Description=license-manager-agent timer
After=network.target

[Timer]
OnActiveSec=15min
Persistent=True

[Install]
WantedBy=timers.target
