# src/templates/license-manager-agent.timer
[Unit]
Description=license-manager-agent timer
After=network.target

[Timer]
OnUnitActiveSec={{ stat_interval }}
Persistent=True

[Install]
WantedBy=timers.target
